# root/lib/udev/rules.d/90-usb-can.rules
ACTION=="add", SUBSYSTEM=="net", ATTRS{idVendor}=="1d50", ATTRS{idProduct}=="606f", RUN+="/usr/bin/env can_add.sh %k"
